# sanger_deployment/backup.logrotate.crontab
52 * * * * /usr/sbin/logrotate /var/opt/citation_reporter_prod/backup.logrotate.conf -s /var/opt/citation_reporter_prod/backup.logrotate.state -f
